# extras-3.3/profile.d/40wine.sh
#!/bin/sh

PATH=$PATH:/opt/wine/bin
MANPATH=$MANPATH:/opt/wine/share/man

export PATH MANPATH
