# extras-3.3/profile.d/40wine.csh
#!/bin/csh

setenv PATH ${PATH}:"/opt/wine/bin"
setenv MANPATH ${MANPATH}:"/opt/wine/share/man"
