# slackware-8.0/profile.d/gnome2.sh
#!/bin/sh

PKG_CONFIG_PATH="$PKG_CONFIG_PATH:/opt/gnome/lib/pkgconfig"
export PKG_CONFIG_PATH
