# ul-3.3/profile.d/04j2sdk.sh
#!/bin/sh

JAVA_HOME=/opt/j2sdk
PATH=$PATH:/opt/j2sdk/bin
MANPATH=$MANPATH:/opt/j2sdk/man

export JAVA_HOME PATH MANPATH
