# ul-3.3/profile.d/10X11.csh
#!/bin/csh

setenv PATH ${PATH}:"/usr/X11R6/bin"
setenv MANPATH ${MANPATH}:"/usr/X11R6/man"
