# ul-3.3/profile.d/10X11.sh
#!/bin/sh

PATH=$PATH:/usr/X11R6/bin
MANPATH=$MANPATH:/usr/X11R6/man

export PATH MANPATH
