# ul-3.3/profile.d/20gnome.csh
#!/bin/csh

setenv PATH ${PATH}:"/opt/gnome/bin"
setenv MANPATH ${MANPATH}:"/opt/gnome/share/man"
setenv INFOPATH ${INFOPATH}:"/opt/gnome/share/info"
setenv PKG_CONFIG_PATH ${PKG_CONFIG_PATH}:"/opt/gnome/lib/pkgconfig"
