# ul-3.3/profile.d/20gnome.sh
#!/bin/sh

PATH=$PATH:/opt/gnome/bin
MANPATH=$MANPATH:/opt/gnome/share/man
INFOPATH=$INFOPATH:/opt/gnome/share/info
PKG_CONFIG_PATH=$PKG_CONFIG_PATH:/opt/gnome/lib/pkgconfig

export PATH MANPATH INFOPATH PKG_CONFIG_PATH
