# ul-3.3/profile.d/21qt.sh
#!/bin/sh

QTDIR=/opt/qt
PATH=$PATH:/opt/qt/bin

export QTDIR PATH
