# ul-3.3/profile.d/21qt.csh
#!/bin/csh

setenv QTDIR /opt/qt
setenv PATH ${PATH}:"/opt/qt/bin"
