# ul-3.3/profile.d/22kde.csh
#!/bin/csh

setenv PATH ${PATH}:"/opt/kde/bin"
setenv MANPATH ${MANPATH}:"/opt/kde/share/man"
