# ul-3.3/profile.d/22kde.sh
#!/bin/sh

PATH=$PATH:/opt/kde/bin
MANPATH=$MANPATH:/opt/kde/share/man

export PATH MANPATH
