# ul-3.3/profile.d/30mozilla.sh
#!/bin/sh

PATH=$PATH:/opt/mozilla/bin
MANPATH=$MANPATH:/opt/mozilla/share/man
PKG_CONFIG_PATH=$PKG_CONFIG_PATH:/opt/mozilla/lib/pkgconfig

export PATH MANPATH PKG_CONFIG_PATH
