# ul-3.3/profile.d/30mozilla.csh
#!/bin/csh

setenv PATH ${PATH}:"/opt/mozilla/bin"
setenv MANPATH ${MANPATH}:"/opt/mozilla/share/man"
setenv PKG_CONFIG_PATH ${PKG_CONFIG_PATH}:"/opt/mozilla/lib/pkgconfig"
